# package/pi5/ai-trap-ble.service
[Unit]
Description=AI Trap BLE GATT Server
Documentation=https://github.com/nature-sense/ai-trap-gstr
# Start after the main trap binary and BlueZ are up
After=network.target bluetooth.service ai-trap.service
Requires=bluetooth.service
BindsTo=ai-trap.service

[Service]
Type=simple
User=root
ExecStart=/opt/ai-trap/ble_gatt.py --name AI-Trap-%i --socket /tmp/ai-trap-ble.sock
Restart=on-failure
RestartSec=10s
# Give the C++ binary time to create the socket
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
